# service/postfix/run
#!/bin/bash
set -e

# Avoid  warning: smtputf8_enable is true, but EAI support is not compiled in
echo "smtputf8_enable = no" >> /etc/postfix/main.cf

# cat >> /etc/postfix/main.cf << EOF
# # limit smtp to loopback interface & compute engine doesn't support ipv6
# inet_interfaces = loopback-only
# inet_protocols = ipv4
# EOF

# Do we want to modify the config first with the script?
# shellcheck disable=SC1091
[ -f /etc/service/postfix/run.config ] && source /etc/service/postfix/run.config

# Transform CONF_ env vars into postfix confs
# E.g. "CONF_MYDOMAIN=example.com" is applied to the main.cf file as as "mydomain=example.com"
printenv | grep '^CONF_' | sed 's/^CONF_\(.*\)=\(.*\)/\L\1\E=\2/' | tr '\n' '\0' | xargs -0 postconf -e

if [[ -n "$MAILNAME" ]]; then
	echo "$MAILNAME" > /etc/mailname
	postconf -e myorigin="/etc/mailname"

	cat >> /etc/postfix/main.cf <<- EOF
	# Force ehlo behavior
	smtp_always_send_ehlo = yes
	smtp_helo_name = $MAILNAME
	EOF
fi

if [[ -n "$ROOT_ALIAS" ]]; then
	if [[ -f /etc/aliases ]]; then
		sed -i '/^root:/d' /etc/aliases
	fi
	echo "root: $ROOT_ALIAS" >> /etc/aliases
	newaliases
fi

if [[ -n "$SASL_AUTH" ]]; then
	cat >> /etc/postfix/main.cf <<- EOF
	smtp_sasl_auth_enable = yes
	smtp_sasl_password_maps = texthash:/etc/postfix/sasl_passwd
	smtp_sasl_security_options = noanonymous
	EOF

	echo "$CONF_RELAYHOST $SASL_AUTH" > /etc/postfix/sasl_passwd
fi

mkdir -p /etc/postfix/pgsql

### Postgres Maps ###
# For each PGQUERY_<NAME>=<sql query> env var, generate /etc/postfix/pgsql/<name>.cf
while read -r line; do
	name="${line%%=*}"
	query="${line#*=}"
	cfname=$(echo "$name" | sed 's/^PGQUERY_//' | tr '[:upper:]' '[:lower:]')
	cat > "/etc/postfix/pgsql/${cfname}.cf" <<- EOF
	hosts = $POSTGRES_HOSTS
	user = $POSTGRES_USER
	password = $POSTGRES_PASSWORD
	dbname = $POSTGRES_ALIAS_DB
	query = $query
	EOF
done < <(printenv | grep '^PGQUERY_')
### END Postgres Maps ###

if [[ -f "/usr/libexec/postfix/master" ]]; then
	cmd="/usr/libexec/postfix/master"
fi

if [[ -f "/usr/lib/postfix/master" ]]; then
	cmd="/usr/lib/postfix/master"
fi

if [[ -z "$cmd" ]]; then
	echo "Could not find postfix master in /usr/lib or /usr/libexec"
	exit 1
fi

# In case of wanting to share one mount for both spool and data (by default in /var/lib/postfix)
mkdir -p /var/spool/postfix/data
chown postfix:root /var/spool/postfix/data
chmod 700 /var/spool/postfix/data

if [[ -n "$POSTCONF_OUTPUT" ]]; then
	postconf >$POSTCONF_OUTPUT
fi

# Wait for rsyslog to be ready, then log startup
while [[ ! -S /dev/log ]]; do sleep 0.1; done
logger -p mail.info "###### Postfix Starting ######"

"$cmd" -c /etc/postfix -d 2>&1
